# packaging/arch/backlayer.install
post_install() {
  cat <<'EOF'
Backlayer installs the UI, daemon, runner binaries, assets, and Plasma wallpaper plugin.

Recommended next step:
  systemctl --user enable --now backlayerd.service

Notes:
  - backlayer-ui can start backlayerd on demand if the daemon socket is missing
  - the user service is still preferred for login autostart and crash recovery
  - KDE Plasma users should select wallpaper type "Backlayer" in Wallpaper Settings
EOF
}

post_upgrade() {
  post_install
}
